# heraldradio

#!/bin/bash
# /etc/init.d/heraldradio

### BEGIN INIT INFO
# Provides:          heraldradio
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: heraldradio playback
# Description:       This service is used to play heraldradio
### END INIT INFO


case "$1" in 
    start)
        echo "Starting heraldradio"
        /home/pi/herald/herald.sh
        ;;
    stop)
        echo "Stopping heraldradio"
        killall herald.sh
        ;;
    *)
        echo "Usage: /etc/init.d/heraldradio start|stop"
        exit 1
        ;;
esac

exit 0
